# target/linux/bcm27xx/base-files/etc/inittab
::sysinit:/lib/systemd/systemd
::shutdown:/lib/systemd/systemd --system --halt
::askconsole:/usr/bin/login
tty1::askfirst:/usr/libexec/login.sh
